# .config/user/dunst.service
[Unit]
Description=Dunst notification daemon

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=default.target

