# .config/user/startpage-update.service
[Unit]
Description=Update startpage system stats and service status

[Service]
Type=oneshot
ExecStart=/bin/bash -c '%h/.local/bin/startpage-stats.sh && %h/.local/bin/startpage-ping.sh'
