# .config/user/startpage-server.service
[Unit]
Description=Startpage local web server
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/startpage
ExecStart=/usr/bin/python -m http.server 7777 --bind 127.0.0.1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
