# system-configs/systemd-system/greetd.service.d/crash-limit.conf
[Unit]
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
# Clear the legacy values from the main unit
StartLimitBurst=
StartLimitInterval=
