# system-configs/systemd-system/greetd.service.d/crash-limit.conf
[Unit]
StartLimitBurst=3
StartLimitIntervalSec=60
