# system-configs/systemd-system/power-switch.service
[Unit]
Description=Apply power profile based on AC/battery state
# powerprofilesctl communicates via D-Bus (socket-activated), no hard ordering needed

[Service]
Type=oneshot
ExecStart=/home/dro/.local/bin/power-switch
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
