# system-configs/udev.rules.d/85-power-switch.rules
# Trigger power-switch on AC plug/unplug
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/home/dro/.local/bin/power-switch"
