# system-configs/systemd-system/logrotate.service.d/no-ac-requirement.conf
[Unit]
# Clear the AC power condition so logrotate runs on battery too
# (Framework 13 laptop — often on battery for extended sessions)
ConditionACPower=
