# system-configs/systemd-system/borg-backup/borg-backup.timer
[Unit]
Description=Borg system backup timer

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
